# recipes-support/carbleurator/carbleurator/carbleurator.service
[Unit]
Description=carbleurator RC car control daemon
After=bluetooth.service
Requires=bluetooth.service
After=hciuart.service
Requires=hciuart.service
After=systemd-udevd.service
Requires=systemd-udevd.service
StartLimitBurst=5
StartLimitIntervalSec=60s
StartLimitAction=poweroff

[Service]
Restart=on-failure
RestartSec=500ms
ExecStart=/usr/bin/carbleurator

[Install]
WantedBy=basic.target
